# ring.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/ring
ExecStop=/bin/bash -c 'rm /run/ringing'

[Install]
WantedBy=multi-user.target
