# twinkle.service
[Service]
ExecStart=/usr/bin/twinkle -c -f /etc/twinkle.cfg
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
