# twinkle.service
[Service]
ExecStart=/usr/bin/twinkle-console /etc/twinkle.cfg
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
